# ollama-copilot.service
[Unit]
Description=Ollama Copilot Proxy
After=network.target

[Service]
WorkingDirectory=%h
ExecStart=%h/.local/bin/ollama-copilot
Restart=on-failure
RestartSec=5
SyslogIdentifier=ollama-copilot

[Install]
# Para serviços de usuário, default.target é o padrão correto
WantedBy=default.target
